# etc/com.expediter.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  LaunchAgent for the Expediter daemon.

  Install:
    1. Build the daemon:           bun run build
    2. Edit paths below if your repo lives anywhere other than
       /Users/asteroidhunter/PycharmProjects/expediter
    3. Symlink (NOT copy) into LaunchAgents so plist edits in the repo
       propagate without re-running install:
         ln -sf /Users/asteroidhunter/PycharmProjects/expediter/etc/com.expediter.daemon.plist \
                ~/Library/LaunchAgents/com.expediter.daemon.plist
    4. Load:  launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.expediter.daemon.plist
    5. First time the user taps a ticket, macOS shows a one-time TCC
       Automation prompt for bun → Terminal.app. Click Allow.

  MUST be a LaunchAgent (gui/<uid> domain), NOT a LaunchDaemon — daemons run
  with no GUI session and cannot drive Terminal.app via osascript, which
  would break tap-to-focus.

  Set ANTHROPIC_API_KEY in ~/.expediter/env (mode 600), not in this plist.
-->
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.expediter.daemon</string>

	<key>ProgramArguments</key>
	<array>
		<string>/Users/asteroidhunter/PycharmProjects/expediter/bin/expediter-daemon.sh</string>
	</array>

	<key>WorkingDirectory</key>
	<string>/Users/asteroidhunter/PycharmProjects/expediter</string>

	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
	</dict>

	<key>RunAtLoad</key>
	<true/>

	<key>KeepAlive</key>
	<true/>

	<key>StandardOutPath</key>
	<string>/Users/asteroidhunter/Library/Logs/expediter.out.log</string>

	<key>StandardErrorPath</key>
	<string>/Users/asteroidhunter/Library/Logs/expediter.err.log</string>

	<key>ProcessType</key>
	<string>Interactive</string>
</dict>
</plist>
